# ansible/roles/authelia/templates/authelia.service.j2
[Unit]
Description=Authelia authentication server
Documentation=https://www.authelia.com
Wants=network-online.target
After=network-online.target

[Service]
User=authelia
Group=authelia
EnvironmentFile=/etc/authelia/authelia.env
ExecStart=/usr/local/bin/authelia --config /etc/authelia/configuration.yml
Restart=on-failure
RestartSec=5s
NoNewPrivileges=true
PrivateTmp=true
ProtectSystem=strict
ReadWritePaths=/var/lib/authelia /var/log/authelia

[Install]
WantedBy=multi-user.target
